# root/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

/usr/sbin/crond -f -d 8 -c /etc/crontabs 2>&1
